# services/activities-full-patch-weekly.timer
[Unit]
Description=Fully patch historical activity reports weekly

[Timer]
# 4am every Sunday local time, after the daily activity and backup jobs.
OnCalendar=Sun 04:00
Persistent=true
AccuracySec=10m
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
