# services/update-files-weekly.timer
[Unit]
Description=Run update-files-weekly with catch-up

[Timer]
# 9am every Sunday local time
OnCalendar=Sun 09:00
# run on next wake if missed
Persistent=true
# group wakeups to save power
AccuracySec=5m
# stagger a bit after wake
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
